# consumer/conf.d/vlc.conf
[program:xterm]
command=su vlcuser -c 'vlc --no-qt-privacy-ask --no-audio --fullscreen %(ENV_RTSP_URL)s vlc://quit'
autorestart=true
startsecs=0
exitcodes=0
