# consumer/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc %(ENV_VNC_PORT)s localhost:5900
autorestart=true
